# services/celery/celery.service
[Unit]
Description=Celery Service
After=network.target
Requires=redis.service

[Service]
Type=forking
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/Jenkins_with_Django
ExecStart=/home/ubuntu/Jenkins_with_Django/env/bin/celery -A company worker --loglevel=info --detach \
          --logfile=/home/ubuntu/Jenkins_with_Django/services/celery/celery.log
ExecStop=/home/ubuntu/Jenkins_with_Django/env/bin/celery -A company worker --loglevel=info --detach --stop
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
